# usbwifi/usbwifi.rules
SUBSYSTEM=="block",ACTION=="add",RUN+="__HOME__/usbwifi-run"
